# deploy/systemd/zk-20-backend.service
[Unit]
Description=ZK-20 古建筑消防预警平台后端服务
After=network.target mosquitto.service

[Service]
Type=simple
User=zk-20
WorkingDirectory=/opt/zk-20/backend
Environment=DATABASE_URL=sqlite:/opt/zk-20/backend/data/zk-20.db
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080
Environment=MQTT_HOST=localhost
Environment=MQTT_PORT=1883
Environment=MQTT_CLIENT_ID=zk-20-backend
Environment=RUST_LOG=info
ExecStart=/opt/zk-20/backend/target/release/zk-20-backend
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=zk-20-backend

[Install]
WantedBy=multi-user.target
